# slopsmith-demucs.service
[Unit]
Description=Slopsmith Demucs Separation Service
After=network.target

[Service]
Type=simple
User=byron
ExecStart=/home/byron/Repositories/slopsmith-demucs-server/.venv/bin/python /home/byron/Repositories/slopsmith-demucs-server/server.py --port 7865
Restart=always
RestartSec=5
Environment=SLOPSMITH_DEMUCS_MODEL=htdemucs

[Install]
WantedBy=multi-user.target
